# config/systemd/user/docker-desktop.service.d/override.conf
[Service]
Environment="DOCKER_DESKTOP_ARGS=--no-window"
